# test/supervisor.conf
[eventlistener:event-listener]
; 默认读取配置文件/etc/supervisor-event-listener.ini
; command=/path/to/supervisor-event-listener
; 指定配置文件路径
; command=/path/to/supervisor-event-listener -c /path/to/supervisor-event-listener.ini

command=/home/cch/mycode/go_projects/src/github.com/ouqiang/supervisor-event-listener/supervisor-event-listener -c /home/cch/mycode/go_projects/src/github.com/ouqiang/supervisor-event-listener/test/supervisor-event-listener.ini
events=PROCESS_LOG,PROCESS_STATE

; for test
directory=/home/cch/mycode/go_projects/src/github.com/ouqiang/supervisor-event-listener/test
priority=1
autostart=true
autorestart=true
; stdout_logfile=/home/cch/mycode/go_projects/src/github.com/ouqiang/supervisor-event-listener/test/listener.log
stderr_logfile=/home/cch/mycode/go_projects/src/github.com/ouqiang/supervisor-event-listener/test/listener.error.log
user=cch



[program:echo]
directory=/home/cch/mycode/go_projects/src/github.com/ouqiang/supervisor-event-listener/test
command=python3 test.py
priority=1
autostart=true
autorestart=true
stdout_logfile=/home/cch/mycode/go_projects/src/github.com/ouqiang/supervisor-event-listener/test/test.log
stderr_logfile=/home/cch/mycode/go_projects/src/github.com/ouqiang/supervisor-event-listener/test/test.error.log
; stdout_events_enabled=true
stderr_events_enabled=true
user=cch


;[group:test]
;programs=echo
